# auditbeat-rules.d/001_self_auditing.conf
#      ___             ___ __      __
#     /   | __  ______/ (_) /_____/ /
#    / /| |/ / / / __  / / __/ __  / 
#   / ___ / /_/ / /_/ / / /_/ /_/ /  
#  /_/  |_\__,_/\__,_/_/\__/\__,_/   
#
# Linux Audit Daemon
# /etc/audit/audit.rules
# 
# Based on rules published here:
#
# https://github.com/bfuzzy/auditd-attack
# https://github.com/Neo23x0/auditd
#

# Self Auditing ---------------------------------------------------------------

## Audit the audit logs
### Successful and unsuccessful attempts to read/write information from the audit records
-w /var/log/auditbeat/ -k auditlog

## Auditd configuration
### Modifications to audit configuration that occur while the audit collection functions are operating
-w /etc/auditbeat/ -p wa -k auditconfig
-w /etc/init.d/auditbeat -p wa -k auditconfig
-w /etc/systemd/system/multi-user.target.wants/auditbeat.service -p wa -k auditconfig

## Monitor for use of audit management tools
-w /bin/auditbeat -p x -k audittools
-w /usr/share/auditbeat/bin/auditbeat -p x -k audittools
-w /usr/share/auditbeat/ -p wa -k auditconfig